# systemd/system/curfew-shutdown.service
[Unit]
Description=Curfew shutdown

[Service]
Type=oneshot
ExecStart=/usr/local/lib/curfew/curfew.sh shutdown
