# systemd/user/curfew-warn.service
[Unit]
Description=Curfew warning notification

[Service]
Type=oneshot
ExecStart=/usr/local/bin/curfew notify
